# crontab
# 練習用：1分ごとに現在時刻をログに出力
* 16 * * * echo "System check: $(date)" >> /var/log/cron.log 2>&1 